# systemd/crc-custom.target
[Unit]
Description=crc custom target
Requires=kubelet-dependencies.target
Requires=crc-env-file-exists.service
After=kubelet-dependencies.target
