# systemd/crc-pullsecret.service
[Unit]
Description=CRC Unit for adding pull secret to cluster
After=crc-wait-apiserver-up.service
StartLimitIntervalSec=450
StartLimitBurst=10

[Service]
Type=oneshot
Restart=on-failure
RestartSec=40
ExecCondition=/usr/local/bin/crc-check-cloud-env.sh
ExecStart=/usr/local/bin/crc-pullsecret.sh

[Install]
WantedBy=multi-user.target
