# systemd/crc-routes-controller.service
[Unit]
Description=CRC Unit starting routes controller
After=crc-env-file-exists.service
After=crc-wait-apiserver-up.service
StartLimitIntervalSec=450
StartLimitBurst=10

[Service]
Type=oneshot
Restart=on-failure
RestartSec=40
ExecCondition=/usr/local/bin/crc-user-mode-networking.sh
ExecCondition=/usr/local/bin/crc-self-sufficient-env.sh
ExecStart=/usr/local/bin/crc-routes-controller.sh

[Install]
WantedBy=crc-custom.target
