# systemd/crc-pullsecret.path
[Unit]
Description=CRC Unit for monitoring the pull secret path
After=kubelet.service

[Path]
PathExists=/opt/crc/pull-secret
TriggerLimitIntervalSec=1min
TriggerLimitBurst=0

[Install]
WantedBy=multi-user.target
