# systemd/ocp-growfs.service
[Unit]
Description=CRC Unit to grow the root filesystem
Requires=crc-custom.target
After=crc-env-file-exists.service

[Service]
Type=oneshot
EnvironmentFile=-/etc/sysconfig/crc-env
ExecCondition=/usr/local/bin/crc-self-sufficient-env.sh
ExecStart=/usr/local/bin/ocp-growfs.sh

[Install]
WantedBy=multi-user.target
