# systemd/ocp-mco-sshkey.service
[Unit]
Description=CRC Unit patching the MachineConfig to add new ssh key
After=crc-wait-apiserver-up.service
StartLimitIntervalSec=1min
StartLimitBurst=1

[Service]
Type=oneshot
Restart=on-failure
ExecCondition=/usr/local/bin/crc-check-cloud-env.sh
ExecStart=/usr/local/bin/ocp-mco-sshkey.sh
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
