# systemd/ocp-userpasswords.path
[Unit]
Description=CRC Unit monitoring the developer and kubeadmin user password file paths
After=kubelet.service

[Path]
PathExists=/opt/crc/pass_kubeadmin
PathExists=/opt/crc/pass_developer
TriggerLimitIntervalSec=1min
TriggerLimitBurst=0

[Install]
WantedBy=multi-user.target
