# systemd/ovs-configuration.service.d/crc-retry-override.conf
[Service]
# CRC: Override retry interval for transient startup races
# Default is 15m which is too long for intermittent failures
# Use 30s for faster recovery without overwhelming the system
Environment="RETRY=30s"
