# systemd/ovs-configuration.service.d/mute-console.conf
[Service]
StandardOutput=append:/var/log/ovs-configure.log
StandardError=append:/var/log/ovs-configure.log
